# docs/deploy/backup-simpleblog.service
[Unit]
Description=Simple Blog DB backup (verified SQLite snapshot)

[Service]
Type=oneshot
ExecStart=/home/brian/scripts/backup-simpleblog.sh
